# digitaldisplay/systemd/digitaldisplay.service
[Unit]
Description=Digital Display Service - Google Slides Kiosk for Raspberry Pi
After=graphical.target network-online.target
Wants=graphical.target network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/node /home/makerspacerp1/digitaldisplay/digitaldisplay/dist/index.js
WorkingDirectory=/home/makerspacerp1/digitaldisplay/digitaldisplay
EnvironmentFile=/home/makerspacerp1/digitaldisplay/digitaldisplay/.env
Restart=always
RestartSec=10
User=makerspacerp1
Environment=NODE_ENV=production
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/makerspacerp1/.Xauthority

[Install]
WantedBy=graphical.target
